# Target/home/audiostream/soundscape/audiostream-ivs-token-refresh.timer
[Unit]
Description=AudioStream - Periodic IVS token refresh

[Timer]
OnBootSec=1min
OnUnitActiveSec=30min
Unit=audiostream-ivs-token-refresh.service
Persistent=true

[Install]
WantedBy=timers.target
