# Target/home/audiostream/soundscape/audiostream.service
[Unit]
Description=AudioStream (umbrella)
After=network-online.target
Wants=network-online.target

# Pull in the actual components.
Wants=audiostream-camilladsp.service
Wants=audiostream-whip.service
Wants=audiostream-ivs-token-refresh.timer

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true

[Install]
WantedBy=multi-user.target
